# image_profiles/gnome/overlay/rootfs/usr/lib/systemd/system/shani-user-setup.path
[Unit]
Description=Watch for new users

[Path]
PathChanged=/etc/passwd

[Install]
WantedBy=multi-user.target
