# image_profiles/shared/overlay/rootfs/etc/systemd/user/org.gnome.GPaste.service.d/override.conf
[Service]
Environment=DISPLAY=:0
Environment=GDK_BACKEND=x11
