# image_profiles/gnome/overlay/rootfs/etc/systemd/system/systemd-zram-setup@.service.d/10-wait-for-swap.conf
[Unit]
After=swap.target

[Service]
ExecStartPre=/usr/abin/wait-for-resume-swap.sh
