# image_profiles/server/overlay/rootfs/etc/profile.d/mcfly-server.sh
# ShaniOS server profile — mcfly smart history defaults
# Sourced by both bash and zsh via /etc/profile.d/
export MCFLY_FUZZY=true
export MCFLY_RESULTS=20
export MCFLY_INTERFACE_VIEW=BOTTOM
export MCFLY_RESULTS_SORT=LAST_RUN
