# image_profiles/server/overlay/rootfs/etc/profile.d/starship-server.sh
# ShaniOS server profile — use server-optimised starship prompt config
# Sourced by both bash and zsh via /etc/profile.d/
export STARSHIP_CONFIG="${STARSHIP_CONFIG:-/etc/starship/server.toml}"
