# image_profiles/server/overlay/rootfs/etc/profile.d/starship-server.sh
# /etc/profile.d/starship-server.sh
# ShaniOS server profile — use server-optimised starship config
export STARSHIP_CONFIG="${STARSHIP_CONFIG:-/etc/starship/server.toml}"
