# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/etc-start-services.service
[Unit]
Description=Start services from overlay filesystem
After=etc-daemon-reload.service

[Service]
Type=oneshot
ExecStart=/usr/abin/start-overlay-services

[Install]
WantedBy=multi-user.target
