# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/etc-daemon-reload.service
[Unit]
Description=Mount overlays and reload systemd daemon
DefaultDependencies=no
After=local-fs.target systemd-remount-fs.service
# beesd-setup.service must complete first so beesd@UUID enable symlink
# exists before the overlay is mounted and daemon-reload picks it up
After=beesd-setup.service
Before=sysinit.target systemd-vconsole-setup.service

[Service]
Type=oneshot
ExecStart=/usr/abin/mount-and-reload
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
