# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/etc-start-services.service
[Unit]
Description=Start services from overlay filesystem
After=etc-daemon-reload.service network.target
Before=display-manager.service
# If stamp already exists, script exits immediately — belt-and-suspenders
ConditionPathExists=!/run/start-overlay-services.done

[Service]
Type=oneshot
ExecStart=/usr/abin/start-overlay-services
RemainAfterExit=yes
TimeoutStartSec=120s
Restart=no

[Install]
WantedBy=multi-user.target
