# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/shani-user-setup.service
[Unit]
Description=Add default groups to newly created users

Requires=accounts-daemon.service etc-overlay.mount
After=accounts-daemon.service etc-overlay.mount

[Service]
Type=oneshot
ExecStart=/usr/abin/shani-user-setup
TimeoutStartSec=30s
