# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/shani-user-setup.path
[Unit]
Description=Watch for new users
After=multi-user.target

[Path]
PathChanged=/etc/passwd
PathChanged=/data/overlay/etc/upper/passwd
TriggerLimitIntervalSec=30s
TriggerLimitBurst=3

[Install]
WantedBy=multi-user.target
