# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/shani-user-setup.path
[Unit]
Description=Watch for new users
After=etc-overlay.mount
Requires=etc-overlay.mount

[Path]
PathChanged=/etc/passwd
Unit=shani-user-setup.service
TriggerLimitIntervalSec=30s
TriggerLimitBurst=3

[Install]
WantedBy=multi-user.target
