# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/snap.mount
# INSTALL TO: /usr/lib/systemd/system/snap.mount
#
# This is a HAND-WRITTEN unit. It must live in /usr/lib/systemd/system/snap.mount
# so it takes precedence over anything fstab-generator produces in
# /run/systemd/generator/snap.mount.
#
# The fstab /snap bind line MUST be removed entirely (not commented out).
# If it remains, fstab-generator regenerates /run/systemd/generator/snap.mount
# on every boot and that generated unit wins over this one.

[Unit]
Description=Bind mount /var/lib/snapd/snap onto /snap
Documentation=man:fstab(5)

# Hard ordering chain:
#   var-lib-snapd.mount                   (@snapd btrfs subvolume)
#     → shanios-tmpfiles-snapd.service    (mkdir /var/lib/snapd/snap + /snap)
#       → snap.mount  (this unit)         (bind /var/lib/snapd/snap → /snap)
#         → snapd.service                 (loop-mounts squashfs into /snap/*/rev/)

After=var-lib-snapd.mount
Requires=var-lib-snapd.mount

After=shanios-tmpfiles-snapd.service
Requires=shanios-tmpfiles-snapd.service

# Ensure the source path is resolvable before we try to bind
RequiresMountsFor=/var/lib/snapd

# nofail: if @snapd is absent on a fresh provision, don't block boot
DefaultDependencies=no

[Mount]
What=/var/lib/snapd/snap
Where=/snap
Type=none
Options=bind

[Install]
WantedBy=snapd.service
WantedBy=local-fs.target
