# image_profiles/shared/overlay/rootfs/usr/lib/systemd/system/snapd.socket.d/override.conf
[Unit]
# Ensure /var/lib/snapd/snap exists inside @snapd subvolume before snapd starts.
# shanios-tmpfiles-snapd.service runs after var-lib-snapd.mount and creates
# the directory persistently inside @snapd.
#
# Ordering:
#   var-lib-snapd.mount
#     → shanios-tmpfiles-snapd.service   (creates /var/lib/snapd/snap in @snapd)
#       → snapd.socket + snapd.service   ← this override
Requires=shanios-tmpfiles-snapd.service
After=shanios-tmpfiles-snapd.service
